# mdi_battery_addon/etc/services.d/battery-service/run
#!/usr/bin/with-contenv bash
set -e

# Log startup message
echo "Starting MDI Battery Add-on..."

# Start the Python application
exec python3 /app/main.py
