# mdi_battery_addon/etc/services.d/battery-service/run
#!/usr/bin/execlineb -P

# Run the Python application
/usr/bin/python3 /app/main.py
